# cvsms-client/cvsms/deploy/supervisor_conf.d/webapp.conf
; =======
;  Django
; =======

[program:webapp]
directory=/var/app/current/
environment=DJANGO_SETTINGS_MODULE="cvsms.settings"
command=gunicorn cvsms.core.wsgi:application -c /var/app/current/deploy/gunicorn.conf
autostart=true
autorestart=true
stdout_logfile = /var/log/webapp/gunicorn_supervisor.log             ; Where to write log messages
redirect_stderr=True
user=webapp